# services/agent-dashboard.service
[Unit]
Description=Multi-Agent Command Center Dashboard
After=network.target

[Service]
Type=simple
User=jarvis
WorkingDirectory=/home/jarvis/dashboard
ExecStart=/home/jarvis/venv/bin/gunicorn --bind 127.0.0.1:5000 --workers 2 --access-logfile - --error-logfile - app:app
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
